# deploy/omniapi.service
[Unit]
Description=OmniAPI
Documentation=https://github.com/cvedix/omniapi
After=network.target

[Service]
Type=simple
User=edgeai
Group=edgeai
WorkingDirectory=/opt/omniapi
ExecStart=/usr/local/bin/omniapi
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
# CPUQuota removed - allow full CPU usage for better performance
# CPUQuota=200%  # Was limiting to 2 cores, removed for better performance

# Security settings
NoNewPrivileges=true
PrivateTmp=true
# Use 'true' instead of 'strict' to allow network binding
ProtectSystem=true
ProtectHome=true
# Allow network access (required for binding ports)
PrivateNetwork=false
ReadWritePaths=/opt/omniapi/logs /opt/omniapi/data /opt/omniapi/uploads /opt/omniapi/instances /opt/omniapi/record /opt/omniapi/models /opt/omniapi/videos /opt/omniapi/fonts /opt/omniapi/run

# Watchdog settings (systemd built-in watchdog)
# NOTE: Disabled because application has its own watchdog implementation
# If you want to use systemd watchdog, implement sd_notify() in the application
# WatchdogSec=30
# NotifyAccess=all

# Environment variables
# Load from .env file if exists, otherwise use defaults
EnvironmentFile=-/opt/omniapi/config/.env
# Config file path - Auto-detection will use /opt/omniapi/config/config.json
# Priority: CONFIG_FILE env var > auto-detection (./config.json > /opt/omniapi/config/config.json)
Environment="API_HOST=0.0.0.0"
Environment="API_PORT=8080"
# Worker executable path (required for subprocess mode)
# Default: edgeos-worker (will search in PATH)
# Production: /usr/local/bin/edgeos-worker
Environment="EDGE_AI_WORKER_PATH=/usr/local/bin/edgeos-worker"
# Socket directory for IPC between main process and workers
# Default: /opt/omniapi/run
Environment="EDGE_AI_SOCKET_DIR=/opt/omniapi/run"
# Execution mode: subprocess (isolated workers) or in-process (legacy)
# Default: subprocess (one worker process per instance - recommended for production)
Environment="EDGE_AI_EXECUTION_MODE=subprocess"
# Instances directory for storing instance configurations
# Default: /opt/omniapi/instances
Environment="INSTANCES_DIR=/opt/omniapi/instances"
# Installation directory (required for finding api-specs/scalar files)
# Default: /opt/omniapi
Environment="OMNIAPI_INSTALL_DIR=/opt/omniapi"
# GStreamer plugin path (required for filesrc and other plugins)
# Priority: 1) Bundled plugins in /opt/omniapi/lib/gstreamer-1.0 (all-in-one)
#          2) Value from /opt/omniapi/config/.env (set by postinst)
#          3) System default path (fallback)
# NOTE: Do NOT hardcode here - let .env file or auto-detection handle it
# The postinst script will set GST_PLUGIN_PATH in .env file automatically
# If .env doesn't exist or doesn't have GST_PLUGIN_PATH, the application will auto-detect
# Disable problematic VA plugins that cause crashes on systems with broken GPU drivers
# Set rank to NONE to prevent GStreamer from loading VA plugins
Environment="GST_PLUGIN_FEATURE_RANK=va:NONE,vaapi:NONE,vaapidecode:NONE,vaapipostproc:NONE"

[Install]
WantedBy=multi-user.target
