# debian/prerm
#!/bin/bash
# prerm script for edgeos-api

set -e

SERVICE_NAME="edgeos-api"

# Stop service before removal
if systemctl is-active --quiet "${SERVICE_NAME}.service" 2>/dev/null; then
    echo "Stopping ${SERVICE_NAME} service..."
    systemctl stop "${SERVICE_NAME}.service" || true
fi

# Disable service
if systemctl is-enabled --quiet "${SERVICE_NAME}.service" 2>/dev/null; then
    echo "Disabling ${SERVICE_NAME} service..."
    systemctl disable "${SERVICE_NAME}.service" || true
fi

exit 0
